# crontabs/dev12.crontab
#! dev12.crontab

# Compress (and delete) log files
14  1   * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3 35
14 7,19 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3

# Reload crontab Monday evening
14 23 * * 1  crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.crontab

# EOF
