# crontabs/optadm.S-01317-C5F.crontab
## optadm.S-01317-C5F.crontab

# Build-info updates for the WiPro Development Environment
#
20 1,3,5,7-17,19,21,23 * * * /cpg/3rdParty/scripts/cpg/buildinfo/updateBuildInfo.sh 2>&1
50 7-17 * * * /cpg/3rdParty/scripts/cpg/buildinfo/updateBuildInfo.sh 2>&1

# Reload crontab Monday evening
10 23 * * 1  crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.`hostname`.crontab

# EOF
