# crontabs/prd20.crontab
#! prd20.crontab

# Compress (and delete) log files
20  1 * * 0  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 5 400

# Reload crontab Monday evening
20 23 * * 1  crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.crontab

# EOF
