# crontabs/stg13.crontab
#! stg13.crontab

# Compress (and delete) log files
16  1 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3 35
16  7 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3
16 19 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3

# Reload crontab Monday evening
16 23 * * 1 /cpg/3rdparty/scripts/cpg/crontabs/$(logname).crontab

# EOF
