# crontabs/stg14.crontab
#! stg14.crontab

# Compress (and delete) log files
18  1 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3 35
18  7 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3
18 19 * * *  /cpg/3rdparty/scripts/cpg/utils/LogManagement_run.sh 3

# Reload crontab Monday evening
18 23 * * 1 /cpg/3rdparty/scripts/cpg/crontabs/$(logname).crontab

# EOF
