# crontabs/stg22.crontab
#! stg22.crontab

# Compress (and delete) log files
24  1 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3 35
24  7 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3
24 19 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3

# Reload crontab Monday evening
24 23 * * 1  crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.crontab

# EOF
