# crontabs/stg21.crontab
#! stg21.crontab

# Compress (and delete) log files
22  1 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3 35
22  7 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3
22 19 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3

# Reload crontab Monday evening
22 23 * * 1 crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.crontab

# EOF
