# crontabs/stg60.crontab
#! stg60.crontab

# Compress (and delete) log files
55  1 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3 35
55  7 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3
55 19 * * *  /cpg/3rdParty/scripts/cpg/utils/logManagement_run.sh 3

# Reload crontab Monday evening
55 23 * * 1 crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.crontab

# EOF
