# crontabs/teamcity.L-03382-C5A.crontab
#!/usr/bin/crontab -u teamcity

# rsync release repo from staging to prd cpodeploy:/cpg/repos
6,16,26,36,46,56 * * * * /cpg/3rdParty/scripts/cpg/utils/migrate_repos

# remove rsync.logs over 7 days old
29 19 * * * /usr/bin/find $HOME/rsync.logs/rsync* -ctime +7 -exec rm {} +

# Reload crontab Monday evening
11 23 * * 1  crontab /cpg/3rdParty/scripts/cpg/crontabs/${LOGNAME}.`hostname`.crontab

# EOF
