# manifests/prd/nodemanager3x.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='nodemanager'>

<service
	name='application/weblogic/nodemanager'
	type='service'
	version='1'>

	<dependency name='multi-user'
		grouping='require_all'
		restart_on='error'
		type='service'>
			<service_fmri value='svc:/milestone/multi-user:default'/>
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh'
		timeout_seconds='30' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='10' >
	</exec_method>

	<exec_method
		type='method'
		name='refresh'
		exec=':kill -3'
		timeout_seconds='10' >
	</exec_method>

	<instance name='prd30' enabled='true'>
		<method_context>
			<method_credential user='prd30' group='devapps' />
		</method_context>
	</instance>

	<instance name='prd31' enabled='true'>
		<method_context>
			<method_credential user='prd31' group='devapps' />
		</method_context>
	</instance>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Multi Instance WebLogic Node Manager Service</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
