# manifests/general/epagentCpodeploy.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='Introscope:epagent'>

<service
	name='application/introscope/epagent'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance/>

	<!-- Need name service resolution and network to connect to WL AdminServer -->
	<dependency
		name='milestones'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/milestone/multi-user-server:default'/>
		<service_fmri value='svc:/milestone/network:default'/>
		<service_fmri value='svc:/milestone/name-services' />
	</dependency>

	<!-- Need the local file system and NFS client to be running -->
	<dependency
		name='file-systems'
		grouping='require_all'
		restart_on='restart'
		type='service'>
		<service_fmri value='svc:/system/filesystem/local:default' />
		<service_fmri value='svc:/network/nfs/client:default' />
		<service_fmri value='svc:/network/nfs/status:default' />
	</dependency>

	<!-- Need network services like DNS for hostname resolution and ldap for user access -->
	<dependency
		name='network'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/dns/client:default' />
		<service_fmri value='svc:/network/ldap/client:default' />
	</dependency>

	<!-- Need crypto-service to be working for SSL operations to work -->
	<dependency
		name='crypto-service'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/cryptosvc' />
	</dependency>

	<method_context>
		<method_credential user='apm' group='devapps' />
	</method_context>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/installs/introscope/9.7.1.0/epagent_cpodeploy/bin/EPACtrl.sh start'
		timeout_seconds='30' >
	</exec_method>

   <exec_method
		type='method'
		name='status'
		exec='/cpg/3rdParty/installs/introscope/9.7.1.0/epagent_cpodeploy/bin/EPACtrl.sh status'
		timeout_seconds='5' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/cpg/3rdParty/installs/introscope/9.7.1.0/epagent_cpodeploy/bin/EPACtrl.sh stop'
		timeout_seconds='30' >
	</exec_method>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Introscope ePAgent</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
