# manifests/general/epagent.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='epagent'>

<service
	name='site/epagent'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<single_instance/>

	<dependency name='multi-user-server'
			grouping='require_all'
			restart_on='error'
			type='service'>
		<service_fmri value='svc:/milestone/multi-user-server:default'/>
	</dependency>

	<dependency name='file-systems'
			grouping='require_all'
			restart_on='error'
			type='path'>
		<service_fmri value='file://localhost/cpg/3rdParty/.mounted'/>
		<service_fmri value='file://localhost/cpg/cpo_var/.mounted'/>
	</dependency>

	<method_context>
		<method_credential user='apm'/>
	</method_context>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/installs/introscope/9.7.1.0/epagent/bin/EPACtrl.sh start'
		timeout_seconds='15' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/cpg/3rdParty/installs/introscope/9.7.1.0/epagent/bin/EPACtrl.sh stop'
		timeout_seconds='10' >
	</exec_method>

   <!-- default service model of 'contract' -->
	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='contract'/>
	</property_group>

	<!--
	  Both action_authorization and value_authorization are needed
	  to allow the framework general/enabled property to be changed
	  when performing action (enable, disable, etc) on the service.
	-->
	<property_group name='general' type='framework'>
		<propval name='value_authorization' type='astring'
				value='solaris.smf.value.onlinechannel' />
		<propval name='action_authorization' type='astring'
				value='solaris.smf.manage.onlinechannel' />
	</property_group>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Introscope ePAgent</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
