# manifests/stg/nodemanager3x.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='nodemanager'>

<service
	name='application/site/nodemanager'
	type='service'
	version='1'>

	<dependency name='multi-user'
		grouping='require_all'
		restart_on='none'
		type='service'>
			<service_fmri value='svc:/milestone/multi-user:default'/>
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh'
		timeout_seconds='30' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh stop'
		timeout_seconds='10' >
	</exec_method>

	<instance name='stg30' enabled='true'>
		<method_context>
			<method_credential user='stg30' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg31' enabled='false'>
		<method_context>
			<method_credential user='stg31' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg32' enabled='true'>
		<method_context>
			<method_credential user='stg32' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg33' enabled='false'>
		<method_context>
			<method_credential user='stg33' group='devapps' />
		</method_context>
	</instance>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient'/>
	</property_group>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Multi Instance WebLogic Node Manager Service</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
