# manifests/stg/nodemanagerAppAddm.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='nodemanager'>

<service
	name='site/nodemanager'
	type='service'
	version='1'>

	<dependency name='multi-user'
		grouping='require_all'
		restart_on='none'
		type='service'>
			<service_fmri value='svc:/milestone/multi-user:default'/>
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh'
		timeout_seconds='30' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh stop'
		timeout_seconds='10' >
	</exec_method>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient'/>
	</property_group>

	<instance name='stg10' enabled='true'>
		<method_context>
			<method_credential user='stg10' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg11' enabled='false'>
		<method_context>
			<method_credential user='stg11' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg12' enabled='true'>
		<method_context>
			<method_credential user='stg12' group='devapps' />
		</method_context>
	</instance>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient'/>
	</property_group>

	<instance name='stg13' enabled='false'>
		<method_context>
			<method_credential user='stg13' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg14' enabled='true'>
		<method_context>
			<method_credential user='stg14' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg20' enabled='true'>
		<method_context>
			<method_credential user='stg20' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg21' enabled='false'>
		<method_context>
			<method_credential user='stg21' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg22' enabled='true'>
		<method_context>
			<method_credential user='stg22' group='devapps' />
		</method_context>
	</instance>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Multi Instance WebLogic Node Manager Service</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
