# manifests/stg/nodemanager5x.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='nodemanager'>

<service
	name='site/nodemanager'
	type='service'
	version='1'>

	<dependency
		name='file-systems'
		grouping='require_all'
		restart_on='refresh'
		type='path'>
		<service_fmri value='file://cpg/content/.mounted' />
		<service_fmri value='file://cpg/cpo_apps/.mounted' />
		<service_fmri value='file://cpg/cpo_var/.mounted' />
		<service_fmri value='file://cpg/3rdParty/.mounted' />
	</dependency>

	<!-- Need name service resolution and network to connect to WL AdminServer -->
	<dependency
		name='milestones'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/milestone/multi-user-server:default'/>
		<service_fmri value='svc:/milestone/network:default'/>
		<service_fmri value='svc:/milestone/name-services' />
	</dependency>

	<!-- Need the local file system and NFS client to be running -->
	<dependency
		name='file-systems'
		grouping='require_all'
		restart_on='restart'
		type='service'>
		<service_fmri value='svc:/system/filesystem/local:default' />
		<service_fmri value='svc:/network/nfs/client:default' />
		<service_fmri value='svc:/network/nfs/status:default' />
	</dependency>

	<!-- Need network services like DNS for hostname resolution and ldap for user access -->
	<dependency
		name='network'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/dns/client:default' />
		<service_fmri value='svc:/network/ldap/client:default' />
	</dependency>

	<!-- Need crypto-service to be working for SSL operations to work -->
	<dependency
		name='crypto-service'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/cryptosvc' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh'
		timeout_seconds='10' >
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/cpg/3rdParty/scripts/cpg/utils/SMFstartNodeManager.sh stop'
		timeout_seconds='5' >
	</exec_method>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient'/>
	</property_group>

	<!--
	  Both action_authorization and value_authorization are needed
	  to allow the framework general/enabled property to be changed
	  when performing action (enable, disable, etc) on the service.
	-->
	<property_group name='general' type='framework'>
		<propval name='value_authorization' type='astring'
				value='solaris.smf.value.onlinechannel' />
		<propval name='action_authorization' type='astring'
				value='solaris.smf.manage.onlinechannel' />
	</property_group>

	<instance name='stg50' enabled='true'>
		<method_context>
			<method_credential user='stg50' group='devapps' />
		</method_context>
	</instance>

	<instance name='stg51' enabled='false'>
		<method_context>
			<method_credential user='stg51' group='devapps' />
		</method_context>
	</instance>
	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>Multi Instance WebLogic Node Manager Service</loctext>
		</common_name>
	</template>

</service>
</service_bundle>
